# dotfiles/myArchBox/crontab
*	*	*	*	6	btrfs blanace start /run/media/haozeke/Storage > /var/log/cron.log